# overlay/etc/update-motd.d/00-sf20
#!/bin/sh
cat <<'EOF'
=====================================================
  Welcome to SourceForge 2.0 - Custom Linux OS
  Minimal. Fast. AI-ready.
=====================================================
EOF
