# os-pack/rootfs/etc/rc.local
#!/bin/sh -e
systemctl enable sf20-provision.service || true
exit 0
