# etc/systemd/system/sf20-firstboot.service
[Unit]
Description=SourceForge 2.0 First Boot (autoclone + boot logo + ASCII MOTD)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sf20/first-boot.sh
RemainAfterExit=yes
ConditionPathExists=!/var/lib/sf20/.firstboot_done

[Install]
WantedBy=multi-user.target
