# overlay/etc/profile.d/10-sf20-main.sh
#!/bin/sh
if [ -f /usr/local/sf20/banner.txt ]; then
  cat /usr/local/sf20/banner.txt
fi
